# bootstrap/nginx.conf
# nginx job file for upstart
# Its strangeness is due to upstart's inability to detect nginx's PID due to forking >2x

description "nginx http daemon"
author "Mark Eirich"

# Listen and start after the vagrant-mounted event
start on vagrant-mounted
stop on runlevel [!2345]

env DAEMON=/usr/sbin/nginx
env PID=/var/run/nginx.pid

pre-start script
	$DAEMON -t
	if [ $? -ne 0 ]
		then exit $?
	fi
	exec $DAEMON
end script

script
	while pidof $DAEMON >/dev/null; do
		sleep 1
	done
end script

post-stop script
	exec start-stop-daemon --stop --pidfile $PID --name nginx --exec $DAEMON --signal QUIT
end script
